# debian/micro-ros-create-ap.prerm
#!/bin/sh

set -x

systemctl stop micro-ros-create-ap.service
systemctl disable micro-ros-create-ap.service