# debian/micro-ros-create-ap.postinst
#!/bin/sh

set -ex

systemctl enable micro-ros-create-ap.service
systemctl start micro-ros-create-ap.service
